# deploy/honeypot-capture.timer
[Unit]
Description=Run the honeypot capture pipeline every 6 hours

[Timer]
OnCalendar=*-*-* 00/6:00:00
Persistent=true
RandomizedDelaySec=180

[Install]
WantedBy=timers.target
